# packaging/debian/lxc-docker.postinst
#!/bin/sh

# Add docker to the init sequence and start it
/usr/sbin/update-rc.d docker defaults
/etc/init.d/docker start
